# nginx-gw/etc/periodic/weekly/geoip-update
#!/bin/sh
#

wkit geoip-update
wkit nginx reload
