# goaccess/etc/periodic/daily/s4-goaccess
#!/bin/sh
#

wkit goaccess restart

### persist ###

. /etc/conf.d/goaccess

logs=`find /var/log/nginx -name "access.log-*[0-9]"`
goaccess $logs --process-and-exit --persist --restore

### archive ###

[ -d $GOACCESS_DATA_DIR/archive ] || mkdir -p $GOACCESS_DATA_DIR/archive

dist=$GOACCESS_DATA_DIR/archive/${logs##*-}
goaccess --output $dist.html --restore
goaccess --output $dist.json --restore
